# etc/services.d/nginx/run
#!/usr/bin/with-contenv sh
echo >&2 "Starting nginx"
nginx -g 'daemon off;'
